# buildroot/board/rockchip/rk3308/fs-overlay/etc/profile.d/RkEnv.sh
export HOME=/userdata
export PATH=$PATH:$HOME:$HOME/bin:/data/bin:/data/bin/rk_pcba_test
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/data/usrlibs

