# buildroot/board/rockchip/rv1126_rv1109/fs-overlay-sysv/etc/profile.d/RkEnv.sh
export HOME=/oem
export PATH=$PATH:$HOME:$HOME/bin:$HOME/usr/bin:$HOME/sbin:$HOME/usr/sbin
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:$HOME/usrlibs:$HOME/usr/lib:$HOME/lib

